# rootfs.template/etc/inittab
::sysinit:/etc/rc.d/rc.init

ttyS0::respawn:/bin/bash
